# fpm/supervisord.conf
[supervisord]
nodaemon = true
logfile = /var/log/supervisord.log
logfile_maxbytes = 10MB
pidfile = /var/run/supervisord.pid

[program:php-fpm]
command = php-fpm7.0
user = root
autostart = true
autorestart = true

[program:memcached]
command = /usr/bin/memcached -p 11211 -u memcache -m 64 -c 1024
user = root
autostart = true
autorestart = true